# tensorboard.conf
[program:tensorboard]
command=tensorboard --logdir /experiments/save/ --bind_all
